# net-misc/sharelist/files/sharelist.init
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="sharelist daemon"
description="Start sharelist during boot"
command=sharelist
command_args=""

depend() {
	need net.*
	after net.*
}

